# deploy/s21notify.service
[Unit]
Description=s21notify (Rust): уведомления Школы 21 в Telegram и MAX
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/s21notify
ExecStart=/opt/s21notify/s21-server
Restart=always
RestartSec=5
Environment=TZ=Europe/Moscow
EnvironmentFile=/opt/s21notify/.env
# порт 80 без root
AmbientCapabilities=CAP_NET_BIND_SERVICE
# LXC всего с 512 МБ
MemoryMax=300M
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
